# packaging/netwatch.service
[Unit]
Description=netwatch – Internet quality evidence collector (24h default)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=@@RUNNER@@
Restart=always
RestartSec=3
Nice=5

[Install]
WantedBy=multi-user.target
